# comparer.service
[Unit]
Description=Create User Automatically Service

[Service]
ExecStart=/home/dina/src/gererUtilisateur/user.sh
